# templates/restic@.timer.j2
[Unit]
Description=Run Restic On a Schedule

[Timer]
OnCalendar={{ systemd_timer_schedule }}

[Install]
WantedBy=timers.target
